# custom-files/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# 强制修正 root 密码
TARGET_HASH='$5$a1grDqnDettfkcMO$27EoNRhxF4vASwsi4xjtQKrzS9bb0yytF6aUDDMtQV7'
exit 0
