# custom-files/etc/uci-defaults/custom-kwrt
#!/bin/sh
# Custom script to run at first boot
# Place your custom commands below

#opkg remove --autoremove luci-app-wizard luci-app-upnp $(opkg list-installed | grep zh-cn | sed -e "s/\s.*$//")
echo '#!/bin/sh /etc/rc.common

# Copyright 2020 Google LLC.
# Copyright (C) 2021 CZ.NIC z.s.p.o. (https://www.nic.cz/)
# SPDX-License-Identifier: Apache-2.0

USE_PROCD=1
START=80

start_service() {
  local state_file
  local port
  local std_err std_out

  config_load tailscale
  config_get_bool std_out "settings" log_stdout 1
  config_get_bool std_err "settings" log_stderr 1
  config_get port "settings" port 41641
  config_get state_file "settings" state_file /etc/tailscale/tailscaled.state
  config_get fw_mode "settings" fw_mode nftables

  /usr/sbin/tailscaled --cleanup

  procd_open_instance
  procd_set_param command /usr/sbin/tailscaled

  # Starting with v1.48.1 ENV variable is required to enable use of iptables / nftables.
  # Use nftables by default - can be changed to 'iptables' in tailscale config
  procd_set_param env TS_DEBUG_FIREWALL_MODE="$fw_mode"

  # Set the port to listen on for incoming VPN packets.
  # Remote nodes will automatically be informed about the new port number,
  # but you might want to configure this in order to set external firewall
  # settings.
  procd_append_param command --port "$port"
  procd_append_param command --state "$state_file"

  procd_set_param respawn
  procd_set_param stdout "$std_out"
  procd_set_param stderr "$std_err"

  procd_close_instance
}

stop_service() {
  /usr/sbin/tailscaled --cleanup
}' > /etc/init.d/tailscale

chmod +x /etc/init.d/tailscale

uci set dhcp.@dnsmasq[0].port='54'
uci set system.@system[0].hostname='HY3000'
uci set network.lan.ipaddr='192.168.6.1'
uci set system.@system[0].version="by å¾®ä¿¡:Mr___zjz/OpenWrt 24.10.4"

uci set wireless.@wifi-iface[0].ssid="Openwrt-2.4G"
uci set wireless.@wifi-iface[0].encryption='psk2'
uci set wireless.@wifi-iface[0].key='password'
# è®¾ç½® 5G WiFi
uci set wireless.@wifi-iface[1].ssid="Openwrt-5G"
uci set wireless.@wifi-iface[1].encryption='psk2'
uci set wireless.@wifi-iface[1].key='password'

# ç¬¬äºŒéƒ¨åˆ†ï¼šè®¾ç½® NPC å®¢æˆ·ç«¯é…ç½®
if [ ! -f /etc/npc-init.flag ]; then
    WAN_IF=$(uci get network.wan.ifname 2>/dev/null || echo "wan")
    WAN_MAC=$(cat /sys/class/net/$WAN_IF/address 2>/dev/null || echo "00:00:00:00:00:00")
    VKEY=${WAN_MAC}

    uci set npc.@npc[0].server_addr="nps.5251314.xyz"
    uci set npc.@npc[0].vkey="$VKEY"
    uci set npc.@npc[0].compress="1"
    uci set npc.@npc[0].crypt="1"
    uci set npc.@npc[0].enable="1"
    uci set npc.@npc[0].server_port="8024"
    uci set npc.@npc[0].protocol="tcp"
    uci commit npc
sed -i 's|conf_Path="/tmp/etc/npc.conf"|conf_Path="/etc/npc.conf"|g' /etc/init.d/npc
# åˆ›å»ºæ–‡ä»¶å¹¶å†™å…¥æ‰€æœ‰å†…å®¹
touch /etc/npc.conf
sed -i "1i[common]" /etc/npc.conf && \
sed -i "2iserver_addr=nps.5251314.xyz:8024" /etc/npc.conf && \
sed -i "3iconn_type=tcp" /etc/npc.conf && \
sed -i "4ivkey=${VKEY}" /etc/npc.conf && \
sed -i "5iauto_reconnection=true" /etc/npc.conf && \
sed -i "6icompress=true" /etc/npc.conf && \
sed -i "7icrypt=true" /etc/npc.conf
    touch /etc/npc-init.flag
    /etc/init.d/npc enable
    /etc/init.d/npc restart
fi

uci set network.lan.ipaddr='192.168.6.1'
uci commit



          # ============================================
          # 7. ä½¿ç”¨ sed ä¿®æ”¹ç‰ˆæœ¬ä¿¡æ¯
          # ============================================
          echo "ðŸ·ï¸ ä½¿ç”¨ sed ä¿®æ”¹ç‰ˆæœ¬ä¿¡æ¯..."
          
          # å®šä¹‰æ–°å€¼ (æ ¹æ®ä½ çš„éœ€æ±‚ä¿®æ”¹è¿™é‡Œ)
          NEW_ID="openwrt "
          NEW_REL="24.10.4 "
          NEW_REV="R$(date +%Y.%m.%d) "
          NEW_DESC=$TARGET_FILE.$NEW_REL.NEW_REV
          
          # --------------------------------------------
          # 1. ä¿®æ”¹ /etc/openwrt_release
          # --------------------------------------------
          # ç¡®ä¿æ–‡ä»¶å­˜åœ¨(ä»Ž ImageBuilder å¤åˆ¶æˆ–ç”Ÿæˆ)
          # æ³¨æ„ï¼šImageBuilder çš„æ–‡ä»¶é€šå¸¸åœ¨ files/etc/openwrt_releaseï¼Œå¦‚æžœä¸å­˜åœ¨éœ€è¦å…ˆåˆ›å»ºæˆ–ä»Žæºç é‡Œæ‰¾
          # è¿™é‡Œå‡è®¾ files/etc/openwrt_release å·²ç»é€šè¿‡ä¹‹å‰çš„æ­¥éª¤å‡†å¤‡å¥½äº†ï¼Œå¦‚æžœæ²¡æœ‰ï¼Œå¯èƒ½éœ€è¦ä»Ž rootfs é‡Œè§£åŽ‹å‡ºæ¥
          # å¦‚æžœæ˜¯ ImageBuilder çŽ¯å¢ƒï¼Œé€šå¸¸å»ºè®®ç›´æŽ¥è¦†ç›–è€Œä¸æ˜¯ sedï¼Œä½†å¦‚æžœä½ åšæŒç”¨ sedï¼š
          
          TARGET_FILE="files/etc/openwrt_release"
          
          # å¦‚æžœæ–‡ä»¶ä¸å­˜åœ¨ï¼Œåˆ›å»ºä¸€ä¸ªåŸºç¡€æ¨¡æ¿ä»¥å… sed æŠ¥é”™
          if [ ! -f "$TARGET_FILE" ]; then
            echo "âš ï¸ $TARGET_FILE ä¸å­˜åœ¨ï¼Œæ­£åœ¨åˆ›å»º..."
            mkdir -p files/etc
            cat > "$TARGET_FILE" <<EOF
DISTRIB_ID='Kwrt'
DISTRIB_RELEASE='24.10-SNAPSHOT'
DISTRIB_REVISION='12.04.2025'
DISTRIB_TARGET='mediatek/filogic'
DISTRIB_ARCH='aarch64_cortex-a53'
DISTRIB_DESCRIPTION='Kwrt 24.10-SNAPSHOT 12.04.2025'
DISTRIB_TAINTS='busybox'
EOF
          fi

          sed -i "s/^DISTRIB_ID=.*/DISTRIB_ID='$NEW_ID'/" "$TARGET_FILE"
          sed -i "s/^DISTRIB_RELEASE=.*/DISTRIB_RELEASE='$NEW_REL'/" "$TARGET_FILE"
          sed -i "s/^DISTRIB_REVISION=.*/DISTRIB_REVISION='$NEW_REV'/" "$TARGET_FILE"
          sed -i "s/^DISTRIB_DESCRIPTION=.*/DISTRIB_DESCRIPTION='$NEW_DESC'/" "$TARGET_FILE"
          
          echo "âœ… $TARGET_FILE ä¿®æ”¹å®Œæˆ"
          
          # --------------------------------------------
          # 2. ä¿®æ”¹ /usr/lib/os-release (OpenWrt ä¸­é€šå¸¸é“¾æŽ¥åˆ° /etc/os-release)
          # --------------------------------------------
          # æ³¨æ„ï¼šImageBuilder ä¸­è¿™ä¸ªæ–‡ä»¶è·¯å¾„é€šå¸¸æ˜¯ usr/lib/os-release
          # å¦‚æžœä¸å­˜åœ¨ï¼ŒåŒæ ·åˆ›å»º
          
          OS_RELEASE_FILE="usr/lib/os-release"
          mkdir -p files/usr/lib
          
          if [ ! -f "$OS_RELEASE_FILE" ]; then
             # å°è¯•å¤åˆ¶ ImageBuilder è‡ªå¸¦çš„ï¼ˆå¦‚æžœæœ‰ï¼‰æˆ–è€…åˆ›å»ºæ–°çš„
             # è¿™é‡Œç®€å•èµ·è§ï¼Œå¦‚æžœä¸å­˜åœ¨å°±åˆ›å»ºä½ æä¾›çš„é‚£ä¸ªæ¨¡æ¿
             cat > "$OS_RELEASE_FILE" <<EOF
NAME="Kwrt"
VERSION="24.10-SNAPSHOT"
ID="kwrt"
ID_LIKE="lede openwrt"
PRETTY_NAME="Kwrt 24.10-SNAPSHOT"
VERSION_ID="24.10-snapshot"
HOME_URL="https://openwrt.ai/"
BUG_URL="https://bugs.openwrt.org/"
SUPPORT_URL="https://forum.openwrt.org/"
BUILD_ID="12.04.2025"
OPENWRT_BOARD="mediatek/filogic"
OPENWRT_ARCH="aarch64_cortex-a53"
OPENWRT_TAINTS="busybox"
OPENWRT_DEVICE_MANUFACTURER="Kiddin'"
OPENWRT_DEVICE_MANUFACTURER_URL="https://openwrt.org/"
OPENWRT_DEVICE_PRODUCT="Generic"
OPENWRT_DEVICE_REVISION="v0"
OPENWRT_RELEASE="Kwrt 24.10-SNAPSHOT 12.04.2025 by Kiddin'"
OPENWRT_BUILD_DATE="1764828800"
EOF
          fi

          # ä½¿ç”¨åŒå¼•å·ä»¥ä¾¿æ‰©å±•å˜é‡ï¼Œåˆ†éš”ç¬¦ç”¨ | é˜²æ­¢å†…å®¹ä¸­æœ‰æ–œæ æŠ¥é”™
          sed -i "s|^NAME=.*|NAME=\"$NEW_ID\"|" "$OS_RELEASE_FILE"
          sed -i "s|^VERSION=.*|VERSION=\"$NEW_REL\"|" "$OS_RELEASE_FILE"
          sed -i "s|^PRETTY_NAME=.*|PRETTY_NAME=\"$NEW_DESC\"|" "$OS_RELEASE_FILE"
          sed -i "s|^VERSION_ID=.*|VERSION_ID=\"$NEW_REL\"|" "$OS_RELEASE_FILE"
          sed -i "s|^BUILD_ID=.*|BUILD_ID=\"$NEW_REV\"|" "$OS_RELEASE_FILE"
          sed -i "s|^OPENWRT_RELEASE=.*|OPENWRT_RELEASE=\"$NEW_DESC by å¾®ä¿¡ï¼šMr___zjz\"|" "$OS_RELEASE_FILE"

          echo "âœ… $OS_RELEASE_FILE ä¿®æ”¹å®Œæˆ"
          
          # éªŒè¯ä¿®æ”¹
          grep -E "DISTRIB_|PRETTY_NAME" files/etc/openwrt_release usr/lib/os-release || true


# Fix aria2 start issue. Refer https://github.com/aria2/aria2/issues/2152#issuecomment-2685168902
sed -i -e 's/section" log/section" log\n        procd_add_jail_mount "\/usr\/lib" #fix "errorCode=1 OSSL_PROVIDER_load 'legacy' failed"/g' /etc/init.d/aria2 
sed -i 's/root:$5$mRqVeJ2Sz8K/CVTM$y6oiLSEeOHx/9ioPe1NKDPVTNa7ubm7wtbXWShElHX0:20427:0:99999:7:::/root:$1$V4UetPzk$CYXluq4wUazHjmCDBCqXF.::0:99999:7:::/g' etc/shadow
sed -i 's/root::20427:0:99999:7:::/root:$1$V4UetPzk$CYXluq4wUazHjmCDBCqXF.::0:99999:7:::/g' etc/shadow

exit 0
